# systemd/scalelite.target
[Unit]
Description=Scalelite
[Install]
WantedBy=multi-user.target